# xe-mailserver/lib/systemd/system/mailserver.target
[Unit]
Description=Dovecot, Exim4 and all other services for a mailserver
StopWhenUnneeded=yes

[Install]
WantedBy=basic.target
