# xe-mailserver/lib/systemd/system/sa-update.service
[Unit]
Description=Update all keys vor spamassassin


[Service]
ExecStart=/bin/systemctl reload spamassassin.service
