# xe-mailserver/lib/systemd/system/sa-update.path
[Unit]

[Path]
PathChanged=/var/lib/spamassassin/compiled/.compiled
